# scripts/osx-package/scripts/launchd/postinstall
#!/bin/bash

set -e
PLIST_ID="com.github.garnierclement.pfe.manticore"

/bin/launchctl load "/Library/LaunchDaemons/$PLIST_ID.plist"