# scripts/osx-package/scripts/launchd/postinstall
#!/bin/bash

set -e
ID="com.github.garnierclement.pfe.manticore"

/bin/launchctl load "/Library/LaunchDaemons/$ID.plist"